# handle-dash-press/handle-dash-press.conf
start on startup

script
    cd /home/matt/src/scripts/handle-dash-press
    handle-dash-press
end script
